# systemd/rcexec.service
[Unit]
Description = irexec

[Service]
ExecStart=/usr/bin/irexec
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target